# ansible/roles/klapp/files/sysdig.service
[Unit]
Description=sysdig service

[Service]
ExecStart=/usr/local/bin/sysdig.init

[Install]
WantedBy=multi-user.target
